# powerctl/files/switch
#!/bin/sh /etc/rc.common
#Automatic startup button power program
START=65

start_service(){
    /usr/bin/switch 
}
